# deploy/supervisor.conf
[program:quantum-backend]
command=/opt/quantum-trading/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8001 --workers 1
directory=/opt/quantum-trading/backend
environment=PATH="/opt/quantum-trading/.venv/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/quantum-backend.err.log
stdout_logfile=/var/log/supervisor/quantum-backend.out.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
startsecs=5
stopwaitsecs=30
